# denyhostspost.cron
59 23 * * * root /ebs/scripts/denyhosts.sh > /dev/null 2>&1
